# LadyDA.conf
[program:LadyDA]
command=python3 /home/curator4/LadyDA/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/LadyDA/test.err.log
stdout_logfile=/home/curator4/LadyDA/test.out.log